# systemd/zenithbot-agent.service.example
[Unit]
Description=ZenithBot agent server
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Zenithbot
Environment=ZENITHBOT_AGENT_DIR=%h/.zenithbot-agent
Environment=ZENITHBOT_AGENT_CWD=%h
Environment=ZENITHBOT_BACKEND=codex
# Set this on real deployments. Do not commit a real token.
# Environment=ZENITHDOCK_AGENT_TOKEN=replace-with-a-long-random-token
ExecStart=%h/Zenithbot/.venv/bin/python %h/Zenithbot/scripts/agent_server.py serve --bind 0.0.0.0 --port 7850
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
